# alpine-base/rootfs/etc/profile.d/banner.sh
#!/usr/bin/env bash

# Auto load the banner file if exist
if [ -f "/banner" ]; then
    cat /banner
fi
